# configs/valhalla.service
[Service]
ExecStart=/usr/bin/node /srv/www/Valhalla.git/server/main.js
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=Valhalla
User=nobody
Group=nobody
Environment=NODE_ENV=production

[Install]
WantedBy=default.target